# roles/kive_server/templates/barman_backup.service.j2
[Unit]
Description=Run barman weekly backup.

[Service]
Environment=MAILTO={{ mail_admins_emails | join(',') }}

User=barman

ExecStart=/bin/barman backup kive
